# templates/varnish.service.erb
[Unit]
Description=Varnish HTTP accelerator
After=syslog.target network.target

[Service]
Type=forking
LimitNOFILE=131072
LimitMEMLOCK=82000
LimitCORE=infinity
PrivateTmp=true

ExecReload=<%= scope['::varnish::vcl_reload'] %>
ExecStart=/usr/sbin/varnishd <%= scope['::varnish::config::jail_opt'] %> \
  -P /var/run/varnish.pid \
  -t 120 \
<% unless [nil, :undefined, :undef, ''].include?(scope['::varnish::instance_name']) -%>
  -n <%= scope['::varnish::instance_name'] %> \
<% end -%>
  -f <%= scope['::varnish::vcl_conf'] %> \
  <%
    if scope['::varnish::version_major'] == 3
      sep = ','
    else
      sep = ' -a '
    end

    if scope['::varnish::listen'].is_a?(Array)
  -%>
  -a <% scope['::varnish::listen'].each_with_index do |addr, idx| %><%= addr %><% unless addr.include?(':') %>:<%= scope['::varnish::listen_port'] %><% end -%><% if idx != (scope['::varnish::listen'].length-1) -%><%= sep %><% end -%><% end %> \
  <% else -%>
  -a <%= scope['::varnish::listen'] %>:<%= scope['::varnish::listen_port'] %> \
  <% end -%>
  -T <%= scope['::varnish::admin_listen'] %>:<%= scope['::varnish::admin_port'] %> \
<% if scope['::varnish::use_tls_proxy'] -%>
  -a proxy=:<%= scope['::varnish::tls_proxy_port'] %>,PROXY \
<% end -%>
<% if scope['::varnish::version_major'] == '3' -%>
  -w <%= scope['::varnish::min_threads'] %>,<%= scope['::varnish::max_threads'] %>,<%= scope['::varnish::thread_timeout'] %> \
<% else -%>
  -p thread_pool_min=<%= scope['::varnish::min_threads'] %> \
  -p thread_pool_max=<%= scope['::varnish::max_threads'] %> \
  -p thread_pool_timeout=<%= scope['::varnish::thread_timeout'] %> \
<% end -%>
<% scope['::varnish::runtime_params'].each do |k,v| -%>
  -p <%= k %>=<%= v %> \
<% end -%>
  -S <%= scope['::varnish::secret_file'] %> \
  -s <%= scope['::varnish::storage_type'] %>,<% if scope['::varnish::storage_type'] == 'file' -%><%= scope['::varnish::storage_file'] %>,<% end -%><%= scope['::varnish::storage_size'] %> \
<% scope['::varnish::storage_additional'].each do |k| -%>
  -s <%= k %> \
<% end -%>

[Install]
WantedBy=multi-user.target
